# etc/profile.d/ssm_10.15_all.csh
setenv SSM_DOMAIN_HOME "/ssm/net/env/master"
