# etc/profile.d/ssm_10.151_all.sh
export SSM_DOMAIN_HOME="/unique/armnssm/ECssm/ssm-domains-base/setup/v_002"
